# petscale.service
[Unit]
Description=Pet Scale Backend Service
After=network.target

[Service]
User=petscale
Group=petscale
WorkingDirectory=/var/www/backend
Environment="PATH=/var/www/backend/venv/bin"
ExecStart=/var/www/backend/venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target 